# files/etc/uci-defaults/99-asu-defaults
uci -q batch << EOI
del wireless.radio0.disabled
del wireless.default_radio0.disabled
del wireless.radio1.disabled
del wireless.default_radio1.disabled
commit wireless
EOI
